# supervisor/conf.d/cointicker.conf
[program:cointicker]
command=/usr/bin/python /opt/coin-ticker/TickerDisplay_Miner.py
autostart=true
autorestart=true
directory=/opt/coin-ticker/
